# deploy/appclimb-backup.service
[Unit]
Description=Back up AppClimb PostgreSQL
Requires=docker.service
After=docker.service

[Service]
Type=oneshot
ExecStart=/opt/apps/appclimb/ops/backup.sh
User=root
